# sources/optware-bootstrap/usbrepo/usbmount.service
[Unit]
Description=Service called when the USB is plugged in

[Service]
Type=oneshot
ExecStart=/usr/bin/optwareUSB add
