# sources/debian-live/config/includes.chroot/etc/profile.d/calnex_livecd.sh
# Copyright © Calnex Solutions, 2014
# All rights reserved.
#
# Define an environment variable to let us know that we are running
# inside a Live CD and that it is expected to get certain failures
# when querying free disk space

export CALNEX_LIVECD=1
