# sources/optware-bootstrap/Springbank/logrotate.cron
#!/bin/sh

test -x /usr/sbin/logrotate || exit 0
su - calnex -c '/usr/sbin/logrotate -s /opt/etc/logrotate.state /opt/etc/logrotate.conf'
